# scripts/launchctl/ai.multica.exchange-monitor-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.multica.exchange-monitor-bot</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/refusned/Desktop/Бот тест Мултика/.venv/bin/python</string>
        <string>/Users/refusned/Desktop/Бот тест Мултика/main.py</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/refusned/Desktop/Бот тест Мултика</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/refusned/Desktop/Бот тест Мултика/logs/launchd.stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/refusned/Desktop/Бот тест Мултика/logs/launchd.stderr.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>SSL_CERT_FILE</key>
        <string>/usr/local/Homebrew/etc/ca-certificates/cert.pem</string>
        <key>REQUESTS_CA_BUNDLE</key>
        <string>/usr/local/Homebrew/etc/ca-certificates/cert.pem</string>
    </dict>
</dict>
</plist>